# orchestrator/systemd/orchestrator.service
[Unit]
Description=Firecracker VM Orchestrator
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/orchestrator
ExecStart=/usr/local/bin/uvicorn orchestrator.main:app --host 0.0.0.0 --port 8090
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
